# voltron.service
# voltron.service — Template systemd unit file
# Copy to /etc/systemd/system/voltron.service and edit paths for your setup.
#
# Replace YOUR_USER, YOUR_GROUP, and paths below to match your deployment:
#   sudo cp voltron.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now voltron

[Unit]
Description=Voltron Agent Dispatcher Worker
After=network.target

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_GROUP
SupplementaryGroups=voltron
WorkingDirectory=/home/YOUR_USER/voltron
ExecStart=/home/YOUR_USER/voltron/.venv/bin/voltron worker
Restart=on-failure
RestartSec=10

# Environment — edit these for your setup
Environment=VOLTRON_BASE_DIR=/home/YOUR_USER/voltron
Environment=VOLTRON_MAX_CONCURRENCY=2
Environment=VOLTRON_AGENT_USER=voltron-agent
Environment=VOLTRON_POLL_INTERVAL=30
Environment=VOLTRON_GITHUB_OWNER=your-github-org
Environment=VOLTRON_CI_CHECK_INTERVAL=60
Environment=VOLTRON_MAX_CI_RETRIES=3
Environment=VOLTRON_ALLOWED_USERS=your-github-username
Environment=VOLTRON_APPROVAL_MODE=review-merge
Environment=VOLTRON_DASHBOARD_PORT=8888
## Set your dashboard password (required to enable dashboard)
# Environment=VOLTRON_DASHBOARD_PASSWORD=changeme

# Uncomment if your repo uses Rust toolchain
# Environment=RUSTUP_HOME=/home/YOUR_USER/.rustup
# Environment=CARGO_HOME=/home/YOUR_USER/.cargo
# Environment=PATH=/home/YOUR_USER/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Security hardening
# NoNewPrivileges omitted — sudo -u voltron-agent needs setuid
ProtectSystem=full
ProtectHome=no
PrivateTmp=yes
PrivateDevices=yes
PrivateIPC=yes
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectClock=yes
RestrictNamespaces=yes
ReadWritePaths=/home/YOUR_USER/voltron/data
ReadWritePaths=/home/YOUR_USER/voltron/repos
ReadWritePaths=/home/YOUR_USER/voltron/logs
ReadWritePaths=/home/voltron-agent

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voltron

[Install]
WantedBy=multi-user.target
